# Backup/crontab
0 * * * * /bin/backup hourly
0 2 * * * /bin/backup daily
0 3 * * 6 /bin/backup weekly
